# scripts/aegis-agent.service
[Unit]
Description=Aegis NGFW Agent v3.0 - System Monitor & Firewall Sync
Documentation=https://github.com/huynhtrungcip/ageis-gnfw
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/aegis
EnvironmentFile=/opt/aegis/.env
ExecStart=/opt/aegis/aegis-agent.sh daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=append:/opt/aegis/agent.log
StandardError=append:/opt/aegis/agent.log
TimeoutStopSec=15

# Security hardening (relaxed — agent needs access to /proc, /sys, docker)
# ProtectSystem=strict causes NAMESPACE errors on many systems
ProtectSystem=false
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
